# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv /bin/bash
set -e

bashio::log.info "Démarrage de nginx sur le port 8099..."

exec nginx -g "daemon off;"
